# java-intellij/idea.conf
[program:idea]
command=/home/experience/idea/bin/idea.sh
autorestart=false
